# parser_admin_bot/parser_admin_bot.service
[Unit]
Description=Crimea Parser Admin Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/parser_admin_bot
EnvironmentFile=/opt/parser_admin_bot/.env
ExecStart=/opt/parser_admin_bot/venv/bin/python /opt/parser_admin_bot/bot.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/parser_admin_bot.log
StandardError=append:/var/log/parser_admin_bot.log

[Install]
WantedBy=multi-user.target
