# _extracted/crimea_parser/deploy/crimea_archive.timer
[Unit]
Description=Crimea Parser — еженедельная архивация result_*.csv

[Timer]
# Понедельник 04:00 MSK — после воскресного прогона парсера и email_finder.
OnCalendar=Mon *-*-* 04:00:00
Persistent=true
Unit=crimea_archive.service

[Install]
WantedBy=timers.target
